# services/otbr/systemd/lisa-otbr-dataset-backup.service
[Unit]
Description=LISA OTBR dataset backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/lisa-edge
ExecStart=/opt/lisa-edge/services/otbr/dataset/backup.sh
